# infra/systemd/obsidian-graph-export.service
[Unit]
Description=Kirobi Obsidian → APK Graph Export
Documentation=https://github.com/Kirobi92/OpenDisruption
After=network.target
Wants=obsidian-graph-export.timer

[Service]
Type=oneshot
User=sven
WorkingDirectory=/home/sven/OpenDisruption
ExecStart=/home/sven/OpenDisruption/infra/scripts/obsidian-graph-export.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120
Nice=10
IOSchedulingClass=idle
# OPE-117: Auto-Sync Obsidian → APK-Graph

[Install]
WantedBy=default.target
